# turtlebot2/app/conf.d/terminator.conf
[program:terminator]
command=terminator
user=developer
autorestart=true
stdout_logfile=/app/log/terminator.log
redirect_stderr=true